# copy/.config/systemd/user/go-to-bed.service
[Unit]
Description=Shuts the computer down to force me to go to bed

[Service]
ExecStart=/usr/bin/shutdown

[Install]
WantedBy=default.target
